# dply.vps-deployster-landing.conf
[program:vps-deployster-landing]
directory=/var/www/vps-deployster-landing
command="/bin/tail -f /dev/null"
autorestart=false
startsecs=0
environment=DPLY_PIPELINE_PORT="80"