# dply.vps-deployster-landing.conf
[program:vps-deployster-landing]
directory=/var/www/vps-deployster-landing
environment=DPLY_PIPELINE_PORT="80"